# cookbooks/graphite/templates/default/carbon-upstart.conf.erb
start on startup
expect daemon

script
cd <%= @node[:graphite][:dir] %>
./bin/carbon-cache.py start
end script

